# miot-pipe.service
[Unit]
Description=Miot MCP Server - Xiaozhi Bridge + HTTP MCP
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi5
WorkingDirectory=/home/pi5/src/miot-skill
EnvironmentFile=/home/pi5/src/miot-skill/.env
ExecStart=/home/pi5/src/miot-skill/venv/bin/python -m miot_skill --xiaozhi --http-port 8300
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
